# claude-slack-bot/supervisord.conf
; supervisord config for the Axon bridge.
;
; Started by the devcontainer postStartCommand on every codespace start:
;   supervisord -c claude-slack-bot/supervisord.conf
; supervisord daemonizes (nodaemon=false) and keeps the bridge alive, restarting
; it on crash with backoff. Central's SSH provisioning runs `supervisorctl
; restart axon-bridge` after it writes the bridge's .env, so the first boot
; (before credentials exist) recovers without waiting for a crash cycle.
;
; All paths are absolute: lifecycle commands and `gh codespace ssh` run from
; different working dirs, so we don't rely on cwd.

[supervisord]
nodaemon=false
logfile=/workspaces/cities/claude-slack-bot/supervisord.log
pidfile=/tmp/coderbots-supervisord.pid
childlogdir=/workspaces/cities/claude-slack-bot

[unix_http_server]
file=/tmp/coderbots-supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/coderbots-supervisor.sock

[program:axon-bridge]
command=/usr/local/bin/uv run python axon_bridge.py
directory=/workspaces/cities/claude-slack-bot
autostart=true
autorestart=true
; Must stay up 10s to count as "started"; a missing/bad .env makes the bridge
; exit immediately, so it retries a few times then goes FATAL (provisioning
; revives it via supervisorctl restart once .env is written).
startsecs=10
startretries=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/workspaces/cities/claude-slack-bot/bridge.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; uv lives in /usr/local/bin; HOME so the bridge's claude subprocess finds the
; user's ~/.claude credentials.
environment=HOME="/home/vscode",PATH="/usr/local/bin:/usr/bin:/bin:/home/vscode/.local/bin"
