# claude-slack-bot/supervisord.conf
; supervisord config for the Ably bridge.
;
; Started by the devcontainer postStartCommand on every codespace start:
;   supervisord -c claude-slack-bot/supervisord.conf
; supervisord daemonizes (nodaemon=false) and keeps the bridge alive, restarting
; it on crash with backoff. On a fresh boot (before Central's provisioning has
; written .env) the bridge doesn't crash — it waits, re-reading .env every few
; seconds until the keys appear (see wait_for_config in bridge.py). Central
; also runs `supervisorctl restart bridge` after writing .env to pick up
; refreshed credentials promptly.
;
; Paths use %(here)s (this config file's own dir), so they resolve correctly
; regardless of the repo name / checkout path and of the caller's cwd.

[supervisord]
nodaemon=false
logfile=%(here)s/supervisord.log
pidfile=/tmp/coderbots-supervisord.pid
childlogdir=%(here)s

[unix_http_server]
file=/tmp/coderbots-supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/coderbots-supervisor.sock

[program:bridge]
command=/usr/local/bin/uv run python bridge.py
directory=%(here)s
autostart=true
autorestart=true
; Must stay up 10s to count as "started". The bridge waits for config rather
; than exiting, so these retries are for genuine crashes, not first-boot.
startsecs=10
startretries=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=%(here)s/bridge.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; uv lives in /usr/local/bin; HOME so the bridge's claude subprocess finds the
; user's ~/.claude credentials.
environment=HOME="/home/vscode",PATH="/usr/local/bin:/usr/bin:/bin:/home/vscode/.local/bin"
